# deployment/macos/io.stevedores.mindofagent.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.stevedores.mindofagent</string>

    <key>ProgramArguments</key>
    <array>
        <string>__MINDOFAGENT_BIN__</string>
    </array>

    <!-- Launch at user login. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on crash, but don't restart on clean exit (so quitting
         from the menu actually quits). -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Only load in the user's Aqua session. Without this key launchd
         also tries to load the agent in LoginWindow / Background, where
         there is no menu bar to render into. -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <!-- Interactive: app touches the user's GUI session, which is the
         whole point of MenuBarExtra. -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <!-- Logs in ~/Library/Logs/MindOfAgent (per-user, no symlink-attack
         surface, survives reboot, picked up by Console.app). LaunchAgent
         plists don't expand ~ or $HOME, so __HOME__ is rendered at
         install time by the Makefile. -->
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/MindOfAgent/mindofagent.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/MindOfAgent/mindofagent.log</string>
</dict>
</plist>
